# deploy/systemd/anicon-backend.service
[Unit]
Description=AniCon Backend (Spring Boot)
After=network.target

[Service]
Type=simple
User=ubuntu

# Env vars loaded from the secrets file (set up once manually on the VM)
EnvironmentFile=/home/ubuntu/anicon-secrets/backend.env

ExecStart=/usr/bin/java -Xms256m -Xmx768m \
    -jar /opt/anicon/anicon-backend.jar \
    --spring.profiles.active=prod
ExecStop=/bin/kill -15 $MAINPID

Restart=on-failure
RestartSec=10
TimeoutStopSec=30

StandardOutput=append:/var/log/anicon/backend.log
StandardError=append:/var/log/anicon/backend-error.log

[Install]
WantedBy=multi-user.target
